# packaging/tinyproxy.service
[Unit]
Description=tinyproxy - secure reverse proxy and web server
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/tinyproxy serve
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tinyproxy

[Install]
WantedBy=multi-user.target
